# crates/daemon/debian/scripts/postinst
#!/bin/sh
set -e

MOOR_USER="moor"
DATABASE_DIR="/var/spool/moor-daemon"
LIB_DIR="/var/lib/moor"
CONFIG_DIR="/etc/moor"
MOOR_HOME=${DATABASE_DIR}

create_user() {
  if id ${MOOR_USER} > /dev/null 2>&1; then return; fi
  adduser --system --home "${MOOR_HOME}" --group ${MOOR_USER}
}

create_directories() {
  mkdir -p ${DATABASE_DIR}
  mkdir -p ${LIB_DIR}
  mkdir -p ${CONFIG_DIR}
  chown -R ${MOOR_USER}:${MOOR_USER} ${LIB_DIR}
  chown -R ${MOOR_USER}:${MOOR_USER} ${DATABASE_DIR}
  chown -R ${MOOR_USER}:${MOOR_USER} ${CONFIG_DIR}
  chmod 700 ${CONFIG_DIR}

  echo ""
  echo "** The moor-daemon package includes a LambdaCore-based core (with modifications for moor)"
  echo "** at /usr/share/moor/cores/lambda-moor/src - see README.md in that directory for details."
  echo "** This core will be used by default. To use a different core, update"
  echo "** /etc/moor/daemon-config.yaml (import_export.input_path) and restart with 'systemctl restart moor-daemon'."
  echo ""
}

generate_keys() {
  # Only generate keys if they don't already exist (e.g., on first install, not upgrades)
  if [ ! -f ${CONFIG_DIR}/moor-signing-key.pem ] || [ ! -f ${CONFIG_DIR}/moor-verifying-key.pem ]; then
    echo "** Generating PASETO signing and verifying keys for client authentication."
    openssl genpkey -algorithm ed25519 -out ${CONFIG_DIR}/moor-signing-key.pem
    openssl pkey -in ${CONFIG_DIR}/moor-signing-key.pem -pubout -out ${CONFIG_DIR}/moor-verifying-key.pem

    # Set permissions on the keys
    chown ${MOOR_USER}:${MOOR_USER} ${CONFIG_DIR}/moor-signing-key.pem
    chown ${MOOR_USER}:${MOOR_USER} ${CONFIG_DIR}/moor-verifying-key.pem
    chmod 600 ${CONFIG_DIR}/moor-signing-key.pem ${CONFIG_DIR}/moor-verifying-key.pem
  else
    echo "** PASETO keys already exist at ${CONFIG_DIR}/moor-*-key.pem, skipping generation."
  fi
}

generate_enrollment_token() {
  # Only generate token if it doesn't already exist
  if [ ! -f ${CONFIG_DIR}/enrollment-token ]; then
    echo "** Generating enrollment token."
    if command -v uuidgen >/dev/null 2>&1; then
      uuidgen > ${CONFIG_DIR}/enrollment-token
    else
      cat /proc/sys/kernel/random/uuid > ${CONFIG_DIR}/enrollment-token
    fi
    chown ${MOOR_USER}:${MOOR_USER} ${CONFIG_DIR}/enrollment-token
    chmod 600 ${CONFIG_DIR}/enrollment-token
  fi
}

case "$1" in
configure)
    create_user
    create_directories
    generate_keys
    generate_enrollment_token
    ;;
esac

## This will only remove masks created by d-s-h on package removal.
deb-systemd-helper unmask moor-daemon.service > /dev/null || true

# was-enabled defaults to true, so new installations run enable.
if deb-systemd-helper --quiet was-enabled moor-daemon.service
then
	# Enables the unit on first installation, creates new
	# symlinks on upgrades if the unit file has changed.
	deb-systemd-helper enable moor-daemon.service > /dev/null || true
	deb-systemd-invoke start moor-daemon
else
	# Update the statefile to add new symlinks (if any), which need to be
	# cleaned up on purge. Also remove old symlinks.
	deb-systemd-helper update-state moor-daemon.service > /dev/null || true
fi

#DEBHELPER#
